# post-installation/aarch64/debian/gnome/resources/system/hci-btaddress@.service
[Unit]
Description=HCI bluetooth address fix

[Service]
Type=simple
ExecStart=/bin/sh -c 'sleep 5 && yes | btmgmt -i %I public-addr F4:A8:0D:BB:ED:78'